# init/cachedproxy.ini
[program:za_cachedproxy]
command=/home/cachedproxy/app/cachedproxy
directory=/home/cachedproxy/app/
process_name=%(program_name)s_%(process_num)02d
user=cachedproxy
group=cachedproxy
numprocs=1
autostart=true
autorestart=true
startsecs=0
stopwaitsecs=10
redirect_stderr=false
stdout_logfile=/var/log/za/cachedproxy/cachedproxy-%(process_num)s.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
stderr_logfile=/var/log/za/cachedproxy/cachedproxy-%(process_num)s.err
stderr_logfile_maxbytes=100MB
stderr_logfile_backups=10